# backend/reminders.cron
* * * * * php /var/www/html/tools/send_reminders.php